# services/ir-keytable.service
[Unit]
Description=Ir-keytable setup

[Service]
Type=oneshot
ExerStart=/usr/bin/ir-keytable -p nec

[Install]
WantedBy=multi-user.target
